# deploy/circuitsight.service
[Unit]
Description=AIT Visual Inspector — Defect Detection Service
After=network.target

[Service]
Type=simple
User=ait
Group=ait
WorkingDirectory=/opt/ait-inspector
ExecStart=/opt/ait-inspector/venv/bin/uvicorn src.api.server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
